# roles/local/templates/neostats/neo-neostats-rotate.service.j2
[Unit]
Description=neostats monthly rotator

[Service]
ExecStart=/opt/neo/neostats/resetneostats.sh
